# docker/root/etc/services.d/shortener/run
#!/usr/bin/execlineb
/usr/bin/with-contenv
url_shortener
